# meshbot.service
[Unit]
Description=Bot Meshtastic-Llama
After=network.target
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=dietpi
Group=dialout
WorkingDirectory=/home/dietpi/bot
ExecStart=/usr/bin/python3 /home/dietpi/bot/main_script.py

# Politique de redémarrage
# - Redémarrer uniquement en cas d'échec (pas pour arrêts intentionnels)
# - Délai de 10 secondes entre les redémarrages
# - Maximum 5 redémarrages en 5 minutes (prévention boucle infinie)
Restart=on-failure
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=300

# Variables d'environnement
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/home/dietpi
Environment=HOME=/home/dietpi

# Gestion des logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=meshtastic-bot

# Sécurité basique
NoNewPrivileges=true
SupplementaryGroups=dialout tty

# Limites de ressources
LimitNOFILE=1024
MemoryMax=512M

[Install]
WantedBy=multi-user.target

